# deploy/pagetally.service
[Unit]
Description=pagetally analytics server
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=pagetally
Group=pagetally
WorkingDirectory=/opt/pagetally
EnvironmentFile=/opt/pagetally/pagetally.env
ExecStart=/opt/pagetally/pagetally-server
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/opt/pagetally
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
